# etc/rc.local
synclient TouchpadOff=1
sudo tlp setcharge 95 100 BAT0
